# deploy/launchd/com.linxz.star-hotel-deal.prewarm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.linxz.star-hotel-deal.prewarm</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.12/bin/python3</string>
    <string>scripts/prewarm_mysql_cache.py</string>
    <string>--limit</string>
    <string>6</string>
    <string>--mode</string>
    <string>full</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/linxiaozhong/development/star-hotel-deal-app</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>StandardOutPath</key>
  <string>/Users/linxiaozhong/development/star-hotel-deal-app/.cache/prewarm.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/linxiaozhong/development/star-hotel-deal-app/.cache/prewarm.err.log</string>
</dict>
</plist>
